# com.transcriptsync.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.transcriptsync.watchdog</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/codyaustin/claude-context-updater/transcriptsync_watchdog.sh</string>
    </array>

    <!-- Run every 30 minutes (1800 seconds) -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Also run immediately when loaded -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Library/Frameworks/Python.framework/Versions/3.11/bin</string>
        <key>HOME</key>
        <string>/Users/codyaustin</string>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/codyaustin/claude-context-updater/logs/launchd_watchdog.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/codyaustin/claude-context-updater/logs/launchd_watchdog_error.log</string>

    <!-- Low priority -->
    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
